# .config/systemd/user/cycle-wallpaper.timer
[Unit]
Description=Cycle to next wallpaper

[Timer]
OnCalendar=*:0/30

[Install]
WantedBy=timers.target
