# .config/systemd/user/cycle-wallpaper.timer
[Unit]
Description=Cycle to next wallpaper

[Timer]
OnCalendar=0/0:30:00

[Install]
WantedBy=timers.target
