# .config/systemd/user/cycle-wallpaper.service
[Unit]
Description=Cycle to next wallpaper
Wants=graphical.target
After=graphical.target

[Service]
Type=oneshot
ExecStart=%h/.scripts/wallpaper.sh
