# .config/systemd/user/cycle-wallpaper.timer
[Unit]
Description=Cycle to next wallpaper

[Timer]
OnCalendar=0/1:00:00
Persistent=false

[Install]
WantedBy=timers.target
