# configs/watch-tablet.service
[Unit]
Description=Starts watch tablet service

[Service]
Restart=always
ExecStart=/usr/local/bin/watch-tablet.sh

[Install]
WantedBy=default.target
